# setup/services/flask_app.service
[Unit]
Description=flask app
After=network.target

[Service]
WorkingDirectory=/home/pi/
ExecStart=/usr/bin/python3 /home/pi/app.py
User=pi

[Install]
WantedBy=multi-user.target

